# deploy/systemd/push-http.service
[Unit]
Description=IM Push System - HTTP API 服务
Documentation=即时消息推送系统 HTTP API（基于 PHP + Swoole）
After=network.target mysql.service redis-server.service
Wants=mysql.service redis-server.service
# 启动失败重试上限（systemd 219 用 StartLimitInterval，230+ 用 StartLimitIntervalSec）
StartLimitBurst=5
StartLimitIntervalSec=30

[Service]
Type=simple
# 项目根目录
WorkingDirectory=/www/push-system/backend
# PHP 可执行文件路径（按实际安装路径调整）
ExecStart=/usr/bin/php /www/push-system/backend/public/index.php
# 优雅停止：发送 SIGTERM 后等待 10 秒
KillSignal=SIGTERM
TimeoutStopSec=10
# 自动重启：崩溃后 3 秒重启
Restart=always
RestartSec=3
# 以指定用户运行（避免 root 权限）
User=www-data
Group=www-data
# 环境变量（可在此覆盖 .env 中的配置）
Environment=APP_ENV=production
# 标准输出重定向到 journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=push-http

[Install]
WantedBy=multi-user.target
