# build/supervisor-build.conf
; ============================================================
; Supervisor 配置 - APK 打包工作进程
; ------------------------------------------------------------
; 部署方式：
;   sudo cp build/supervisor-build.conf /etc/supervisor/conf.d/build-worker.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start build-worker
; ============================================================

[program:build-worker]
; 进程名
process_name=%(program_name)s
; 启动命令：以 PHP CLI 常驻运行 BuildWorker
command=/usr/bin/php %(here)s/../build/queue/BuildWorker.php
; 工作目录为 build/queue
directory=%(here)s/../build/queue
; 以哪个用户运行（建议非 root，按需修改）
user=www-data
; 进程数（打包为串行任务，保持 1 即可）
numprocs=1
; 随 supervisor 自动启动
autostart=true
; 异常退出自动重启
autorestart=true
; 启动重试次数
startretries=5
; 退出码白名单（仅这些码之外才视为失败重启）
exitcodes=0,2
; 重启等待秒数
startsecs=3
; 停止等待秒数
stopwaitsecs=30
; 停止信号：SIGTERM（worker 会优雅退出）
stopsignal=TERM

; 标准输出日志
stdout_logfile=%(here)s/../build/queue/worker.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=5
; 标准错误日志
stderr_logfile=%(here)s/../build/queue/worker.err.log
stderr_logfile_maxbytes=20MB
stderr_logfile_backups=5
; 重启时不截断日志
redirect_stderr=false
